# smartpalms-kiosk.service
[Unit]
Description=Smart Palms Kiosk Application
After=network.target graphical.target

[Service]
Type=simple
User=smartpalms
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/smartpalms/.Xauthority
WorkingDirectory=/home/smartpalms/smartpalms-kiosk
ExecStart=/bin/bash /home/smartpalms/smartpalms-kiosk/run_kiosk.sh
Restart=always
RestartSec=10

[Install]
WantedBy=graphical.target 
